# distros/artix-openrc/grow-rootfs.openrc
#!/sbin/openrc-run

description="Grow root filesystem to fill disk (first boot)"

depend() {
    need localmount
}

start() {
    ebegin "Growing root filesystem"
    /usr/local/sbin/grow-rootfs
    eend $?
}
